# scripts/com.backburner.mexc-cookie.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.backburner.mexc-cookie</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npx</string>
        <string>tsx</string>
        <string>/Users/roti/.claude-worktrees/Backburner/nifty-lewin/scripts/mexc-cookie-daemon.ts</string>
        <string>--once</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/roti/.claude-worktrees/Backburner/nifty-lewin</string>

    <!-- Run every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Also run at load -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/roti/.claude-worktrees/Backburner/nifty-lewin/.mexc-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/roti/.claude-worktrees/Backburner/nifty-lewin/.mexc-daemon.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
